# deployment/adr-app.service
[Unit]
Description=Architecture Decisions Application
Documentation=https://decisionrecords.org
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

# Environment file with database URL and other config
EnvironmentFile=/etc/adr-app/env

# Pull latest image before starting (ensures we have newest version)
ExecStartPre=/usr/bin/docker pull adrregistry2024eu.azurecr.io/architecture-decisions:latest

# Remove any existing container (ignore errors if it doesn't exist)
ExecStartPre=-/usr/bin/docker rm -f adr-app

# Start container with environment variables
# Port 8000 exposed for Application Gateway
ExecStart=/usr/bin/docker run --rm --name adr-app \
    -p 8000:8000 \
    --env-file /etc/adr-app/env \
    adrregistry2024eu.azurecr.io/architecture-decisions:latest

# Stop container gracefully
ExecStop=/usr/bin/docker stop -t 30 adr-app

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=adr-app

[Install]
WantedBy=multi-user.target
